# recipes-core/udev/udev-rules-imx/backlight.rules
KERNEL=="backlight", RUN+="/etc/udev/rules.d/set-backlight.sh"
